# docker/etc/s6-overlay/s6-rc.d/04-nqptp/run
#!/command/with-contenv sh
echo "Starting nqptp"
if [ -z ${IPV4ONLY+x} ]; then
    exec /usr/local/bin/nqptp
else
    exec /usr/local/bin/nqptp -ipv4-only
fi
