# docker/etc/s6-overlay/s6-rc.d/02-avahi/run
#!/bin/sh
avahi-daemon --check
if [ $? -eq 0 ]; then
  echo "Reloading avahi-daemon."
  exec avahi-daemon --no-chroot --reload
else
  echo "Starting avahi-daemon."
  exec avahi-daemon --no-chroot
fi