# docker/etc/services.d/02-avahi/run
#!/bin/sh
avahi-daemon --check
if [ $? -eq 0 ]; then
  echo "Reloading avahi-daemon."
  avahi-daemon --no-chroot --reload
else
  echo "Starting avahi-daemon."
  avahi-daemon --no-chroot
fi