# docker/etc/services.d/01-dbus-daemon/finish
#!/usr/bin/execlineb -P
rm -f /var/run/dbus.pid