# docker/etc/services.d/01-dbus-daemon/run
#!/usr/bin/execlineb -P
dbus-daemon --system --nofork