# docker/etc/s6-overlay/s6-rc.d/01-dbus/run
#!/command/with-contenv sh
rm -rf /var/run/dbus.pid
echo "Starting dbus"
exec dbus-daemon --system --nofork